# hddtemp/hddtemp.service
[Unit]
Description=Hard drive temperature monitor daemon
Documentation=man:hddtemp(8)

[Service]
Environment=HDDTEMP_OPTIONS=--listen=127.0.0.1
EnvironmentFile=-/etc/default/hddtemp
ExecStart=/usr/sbin/hddtemp -d -F $HDDTEMP_OPTIONS

[Install]
WantedBy=multi-user.target
