# qbittorrent/qbittorrent.service
[Unit]
Description=qBitorrent client
After=syslog.target network.target local-fs.target
Wants=local-fs.target

[Service]
EnvironmentFile=/etc/default/qbittorrent
User=qbittorrent
Group=qbittorrent
ExecStart=/usr/bin/qbittorrent-nox --webui-port=$DAEMON_PORT

[Install]
WantedBy=multi-user.target
