# qbittorrent/qbittorrent.service
[Unit]
Description=qBitorrent client
After=syslog.target network.target local-fs.target
Wants=local-fs.target
ConditionPathExists=/etc/rc.d/rc.qbittorrent

[Service]
ExecStop=/etc/rc.d/rc.qbittorrent stop
ExecStart=/etc/rc.d/rc.qbittorrent start
ExecReload=/etc/rc.d/rc.qbittorrent restart
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
