# clamav/clamav-update.cron
#!/bin/sh
# Cron script to update ClamAV. Should be located in cron.daily.

(
	# be nice to server and wait a random while so that not
	# all ClamAV installations pound the server at once.

	sleep $((RANDOM%590+10))

	#http_proxy=http://www-cache:3128
	/usr/bin/freshclam --quiet
) &
# relegate to background so the waiting doesn't hog other cron jobs.
