# openssh/systemd/sshd@.service
[Unit]
Description=OpenSSH per-connection server daemon.
After=syslog.target

[Service]
ExecStart=-/usr/sbin/sshd -i
StandardInput=socket
