# samba/systemd/winbindd.service
[Unit]
Description=Samba Winbind daemon
After=syslog.target network.target

[Service]
ExecStart=/usr/sbin/winbindd -F
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
