# mariadb/systemd/mysql-compat.service.in
[Unit]
Description=MySQL compatibility service (another name for @DAEMON_NAME@.service; you should use @DAEMON_NAME@.service instead)
BindsTo=@DAEMON_NAME@.service

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
Also=@DAEMON_NAME@.service
