# brofficeorg/profile.d/brofficeorg.csh
#!/bin/csh
# Force interface to use: gen, gtk or kde
# if ( ! $?SAL_USE_VCLPLUGIN ) then
#     setenv SAL_USE_VCLPLUGIN gtk
# endif
# if ( ! $?OOO_FORCE_DESKTO ) then
#     setenv OOO_FORCE_DESKTOP gnome
# endif
