# gksu-polkit/gksu.service
[Unit]
Description=Gksu server

[Service]
Type=dbus
BusName=org.freedesktop.Gksu
ExecStart=/usr/sbin/gksu-server
