# nfs-utils/systemd/nfs.target
[Unit]
Description=Network File System Server
Requires=var-lib-nfs-rpc_pipefs.mount proc-fs-nfsd.mount rpcbind.service
After=network.target named.service 

[Install]
WantedBy=multi-user.target
