# gdk-pixbuf2/update-gdk-pixbuf-loaders.service
[Unit]
Description=Run update-gdk-pixbuf-loaders
Before=display-manager.service
After=ldconfig.service gio-modules.service glib-compile-schemas.service

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/update-gdk-pixbuf-loaders
StandardOutput=null

[Install]
WantedBy=graphical.target
