# mpd/mpd.service
[Unit]
Description=Music Player Daemon
After=syslog.target network.target sound.target local-fs.target
Requires=sound.target

[Service]
ExecPreStart=-/usr/bin/mpd-start-user-config
ExecStart=/usr/bin/mpd --no-daemon
User=mpd
Group=mpd

[Install]
WantedBy=multi-user.target
