# transmission/transmission.service
[Unit]
Description=Transmission BitTorrent client
After=network.target local-fs.target

[Service]
ExecStop=/etc/rc.d/rc.transmission stop
ExecStart=/etc/rc.d/rc.transmission start
ExecReload=/etc/rc.d/rc.transmission restart
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
